# server/supervisor/conf.d/redis.conf
[program:redis]
directory=server
command=redis-server
autostart=true
autorestart=false
redirect_stderr=true
stopasgroup=true
